# computers/raspi5/root/etc/systemd/system/mount-nas.timer
[Unit]
Description=Mount NAS Timer

[Timer]
OnBootSec=20s
OnUnitActiveSec=5m
Unit=mount-nas.service

[Install]
WantedBy=timers.target
