# computers/raspi5/root/etc/systemd/system/mount-nas.service
[Unit]
Description=Mount NAS

[Service]
Type=oneshot
ExecStart=/home/romain/.local/bin/mount-nas
User=root
